# cookbooks/zookeeperbox/templates/default/zoo.upstart.conf.erb
#!upstart
description "Zookeeper Server"

env USER=vagrant

start on startup
stop on shutdown

respawn

exec sudo -u $USER sh -c "/home/vagrant/zookeeper/bin/zkServer.sh start-foreground"
